# hotshots@.service
[Unit]
Description=HotShots %i scraper

[Service]
Type=oneshot
StandardOutput=journal
ExecStart=/path/to/binary %i

[Install]
WantedBy=default.target